# etc/systemd/system/getty@.service.d/debugmode.conf
[Unit]
ConditionKernelCommandLine=debug

[Service]
ExecStartPre=-/bin/echo 'Debug mode autologin'
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I 38400 linux
